# rover/jetson_stable/stable/services/jetson-mavros.service
[Unit]
Description=MAVROS2 - MAVLink ROS 2 Bridge for Cube Orange
After=network.target

[Service]
Type=simple
User=jay
WorkingDirectory=/home/jay
Environment="ROS_DOMAIN_ID=42"
ExecStart=/bin/bash -c "source /opt/ros/humble/setup.bash && ros2 launch mavros apm.launch fcu_url:=/dev/ttyACM0:57600"
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/jetson-mavros.log
StandardError=append:/var/log/jetson-mavros.log

[Install]
WantedBy=multi-user.target
